# launchd/com.example.stay-awake.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.example.stay-awake</string>
  <key>ProgramArguments</key>
  <array>
    <string>@HOME@/bin/stay-awake.sh</string>
    <string>forever</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/stay-awake.launchd.out</string>
  <key>StandardErrorPath</key>
  <string>/tmp/stay-awake.launchd.err</string>
</dict>
</plist>
